# Anemone/Test Data/info/cydia.list
/.
/Library
/Library/LaunchDaemons
/Library/LaunchDaemons/com.saurik.Cydia.Startup.plist
/etc
/etc/apt
/etc/apt/sources.list.d
/etc/apt/trusted.gpg.d
/usr
/usr/libexec
/usr/libexec/cydia
/usr/libexec/cydia/cfversion
/usr/libexec/cydia/cydo
/usr/libexec/cydia/du
/usr/libexec/cydia/finish.sh
/usr/libexec/cydia/firmware.sh
/usr/libexec/cydia/free.sh
/usr/libexec/cydia/move.sh
/usr/libexec/cydia/setnsfpn
/usr/libexec/cydia/startup
/var
/var/lib
/var/lib/cydia
