# Anemone/Test Data/info/openssh.list
/.
/Library
/Library/LaunchDaemons
/Library/LaunchDaemons/com.openssh.sshd.plist
/etc
/etc/ssh
/etc/ssh/moduli
/etc/ssh/ssh_config
/etc/ssh/sshd_config
/usr
/usr/bin
/usr/bin/scp
/usr/bin/sftp
/usr/bin/ssh
/usr/bin/ssh-add
/usr/bin/ssh-agent
/usr/bin/ssh-keygen
/usr/bin/ssh-keyscan
/usr/libexec
/usr/libexec/sftp-server
/usr/libexec/ssh-keysign
/usr/libexec/ssh-pkcs11-helper
/usr/sbin
/usr/sbin/sshd
/var
/var/empty
